# .crontab
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO="kazuhiko.hrmt@gmail.com"

# そらぴよ
# * * * * * for i in $(seq 5 15 59); do (sleep $i; $HOME/bin/cron/sorapiyo_notify.sh) & done

# 超!A&G+ ################################ 

# Health Check (毎日06:10)
10 06 * * * $HOME/bin/cron/record_agqr_test.sh home log

# 徳井青空のまぁるくなぁれ！ (毎週木曜25:00-25:30)
00 01 * * 5 $HOME/bin/cron/record_agqr_maru.sh home log

# HiBiKi Radio Station ###################

# 徳井青空のまぁるくなぁれ！ (毎週金曜12:00頃)
00 13 * * 5 $HOME/bin/cron/download_hibikiradio_marunare.sh home log
